# sample/openvpn.ini
[program:openvpn]
command=openvpn --config /etc/openvpn/openvpn.conf
autorestart=true
user=root
stdout_events_enabled = true
stderr_events_enabled = true
startsecs=0
